# backend/billmint-backend.service
[Unit]
Description=BillMint Backend API Server
After=network.target

[Service]
Type=simple
User=billmint
WorkingDirectory=/opt/billmint-backend
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="DB_PATH=/var/lib/billmint/db.sqlite"
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=billmint-backend

[Install]
WantedBy=multi-user.target
